# projects/ed-autojump/scripts/ed-afk-nightly.xml
<?xml version="1.0" encoding="UTF-8"?>
<!--
  Task Scheduler XML for the ed-afk nightly runner.

  This is a MANUAL-IMPORT template — review the paths below before
  registering. Two paths likely need editing for your machine:

    1. <Arguments> — the script path. Default points to
       G:\Documents\GIT\ED-AFK\projects\ed-autojump\scripts\nightly-run.ps1
    2. <WorkingDirectory> — the project root.

  Register with:
    Register-ScheduledTask -Xml (Get-Content -Path .\ed-afk-nightly.xml -Raw) `
        -TaskName "ED-AFK Nightly" -User $env:USERNAME

  Or via the Task Scheduler GUI: Action -> Import Task -> select this file.

  The trigger fires daily at 23:00 local time. The 8-hour
  ExecutionTimeLimit is a hard kill safety net; the bot's own
  --duration is the soft limit.
-->
<Task version="1.4" xmlns="http://schemas.microsoft.com/windows/2004/02/mit/task">
  <RegistrationInfo>
    <Author>ed-afk</Author>
    <Description>ED-AFK overnight unattended ed-autojump session.</Description>
    <URI>\ED-AFK Nightly</URI>
  </RegistrationInfo>
  <Triggers>
    <CalendarTrigger>
      <StartBoundary>2026-05-22T23:00:00</StartBoundary>
      <Enabled>true</Enabled>
      <ScheduleByDay>
        <DaysInterval>1</DaysInterval>
      </ScheduleByDay>
    </CalendarTrigger>
  </Triggers>
  <Principals>
    <Principal id="Author">
      <LogonType>InteractiveToken</LogonType>
      <RunLevel>LeastPrivilege</RunLevel>
    </Principal>
  </Principals>
  <Settings>
    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>
    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>
    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>
    <AllowHardTerminate>true</AllowHardTerminate>
    <StartWhenAvailable>false</StartWhenAvailable>
    <RunOnlyIfNetworkAvailable>true</RunOnlyIfNetworkAvailable>
    <IdleSettings>
      <StopOnIdleEnd>false</StopOnIdleEnd>
      <RestartOnIdle>false</RestartOnIdle>
    </IdleSettings>
    <AllowStartOnDemand>true</AllowStartOnDemand>
    <Enabled>true</Enabled>
    <Hidden>false</Hidden>
    <RunOnlyIfIdle>false</RunOnlyIfIdle>
    <WakeToRun>false</WakeToRun>
    <ExecutionTimeLimit>PT8H</ExecutionTimeLimit>
    <Priority>7</Priority>
  </Settings>
  <Actions Context="Author">
    <Exec>
      <Command>powershell.exe</Command>
      <Arguments>-NoProfile -ExecutionPolicy Bypass -File "G:\Documents\GIT\ED-AFK\projects\ed-autojump\scripts\nightly-run.ps1" -DurationHours 6</Arguments>
      <WorkingDirectory>G:\Documents\GIT\ED-AFK\projects\ed-autojump</WorkingDirectory>
    </Exec>
  </Actions>
</Task>
